# config/crontab/crontab
0,30 * * * * ~/work/themis/scripts/crawl_feeds
